# deploy/radiator-api.service
[Unit]
Description=Radiator API
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=exec
User=radiator
Group=radiator
WorkingDirectory=/opt/radiator-api
Environment=PATH=/opt/radiator-api/.venv/bin
ExecStart=/opt/radiator-api/.venv/bin/uvicorn radiator.main:app --host 0.0.0.0 --port 8000
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=radiator-api

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/radiator-api/logs /opt/radiator-api/uploads

[Install]
WantedBy=multi-user.target
